# deploy/mathgame-web.service
[Unit]
Description=Math Game Flask Web App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/mathgame
Environment="PATH=/opt/mathgame/venv/bin"
EnvironmentFile=/opt/mathgame/.env
ExecStart=/opt/mathgame/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
